# data/butter-schedule-snapshot.timer
[Unit]
Description=Butter Schedule Snapshot

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
