# data/butter-schedule-prune.timer
[Unit]
Description=Butter Schedule Prune

[Timer]
OnBootSec=15m
OnUnitActiveSec=8h
RandomizedDelaySec=30m

[Install]
WantedBy=timers.target
