# data/butter-schedule-prune.service.in
[Unit]
Description=Butter Schedule Prune

[Service]
Type=simple
ExecStart=@BINDIR@/butter schedule prune
