# data/butter-schedule-snapshot.service.in
[Unit]
Description=Butter Schedule Snapshot

[Service]
Type=simple
ExecStart=@BINDIR@/butter schedule snapshot
